# data/systemd/cairo-dock.service.in
[Unit]
Description=Cairo-Dock
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cairo-dock
BusName=org.cairodock.CairoDock

[Install]
WantedBy=graphical-session.target
